# deploy/pi-ai-hat-web.service
[Unit]
Description=pi-ai-hat Web Server
After=network.target hailo-daemon.service pi-ai-hat-hotspot.service
Wants=hailo-daemon.service pi-ai-hat-hotspot.service

[Service]
Type=simple
WorkingDirectory=/opt/pi-ai-hat
# Block until the daemon is actually listening — prevents the race where the
# web server starts before hailo_daemon.py has bound to port 8001.
ExecStartPre=/bin/bash -c 'i=0; until ss -tlnp | grep -q ":8001"; do i=$((i+1)); [ $i -ge 60 ] && exit 1; sleep 1; done'
ExecStart=/opt/pi-ai-hat/web_env/bin/python /opt/pi-ai-hat/server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
